# .local/runit/sv/iptables-persistent/run
#!/bin/sh

exec 2>&1; set -e

[ -r conf ] && . ./conf

exec sh -c "/usr/bin/iptables-restore $PATH_TO_RULES_V4; /usr/bin/ip6tables-restore $PATH_TO_RULES_V6; pause"
